# local_docker_dev/conf.d/supervisord.conf
; sudo supervisorctl restart nginx
[program:nginx]
command=nginx -g 'daemon off;'
; stdout_logfile=/dev/stdout
stdout_logfile=/var/log/supervisor/nginx.stdout.log
stdout_logfile_maxbytes=0
; stderr_logfile=/dev/stderr
stderr_logfile=/var/log/supervisor/nginx.stderr.log
stderr_logfile_maxbytes=0
autorestart=true
startretries=10
stopasgroup=true
killasgroup=true

; Disable start php-fpm in local_dev_docker/nginx/entrypoint.sh
; ; for adminer481.php in nginx
; [program:php-fpm]
; command=/usr/sbin/php-fpm8.3 --nodaemonize
; autostart=true
; autorestart=true

# ; sudo supervisorctl restart laravel
# [program:laravel]
# command=php artisan serve --host=0.0.0.0 --port=8000
# directory=/home/ubuntu/code/demo/backend
# ; stdout_logfile=/dev/stdout
# stdout_logfile=/var/log/supervisor/laravel.stdout.log
# stdout_logfile_maxbytes=0
# ; stderr_logfile=/dev/stderr
# stderr_logfile=/var/log/supervisor/laravel.stderr.log
# stderr_logfile_maxbytes=0
# autorestart=true
# startretries=10
# stopasgroup=true
# killasgroup=true

# ; ## Disable schedule in laravel 7 is not support.
# ; sudo supervisorctl restart schedule
# ; [program:schedule]
# ; command=php artisan schedule:work
# ; directory=/home/ubuntu/code/demo/backend
# ; stdout_logfile=/var/log/supervisor/schedule.stdout.log
# ; stdout_logfile_maxbytes=0
# ; stderr_logfile=/var/log/supervisor/schedule.stderr.log
# ; stderr_logfile_maxbytes=0
# ; autorestart=true
# ; startretries=5
# ; stopasgroup=true
# ; killasgroup=true

# ; sudo supervisorctl start "queue:*"
# ; or in dev. to command, artisan queue:listen
# [program:queue]
# process_name=%(program_name)s_%(process_num)02d
# command=php artisan queue:work -v --sleep=3 --tries=3 --max-time=3600
# ; command=php artisan queue:listen -v --tries=1
# directory=/home/ubuntu/code/demo/backend
# stdout_logfile=/var/log/supervisor/queue.stdout.log
# stdout_logfile_maxbytes=0
# stderr_logfile=/var/log/supervisor/queue.stderr.log
# stderr_logfile_maxbytes=0
# autostart=true
# autorestart=true
# startretries=10
# stopasgroup=true
# killasgroup=true
# numprocs=1

# ; Please execute the command, "php artisan optimize" after Updated ".env" file.
# ; sudo supervisorctl start "consumeDTQ:*"
# ; DTQ = DriverTicketsQueue
# [program:consumeDTQ]
# process_name=%(program_name)s_%(process_num)02d
# command=php artisan consume:DriverTicketsQueue
# directory=/home/ubuntu/code/demo/backend
# stdout_logfile=/var/log/supervisor/consumeDTQ_%(process_num)02d.stdout.log
# stdout_logfile_maxbytes=0
# stderr_logfile=/var/log/supervisor/consumeDTQ_%(process_num)02d.stderr.log
# stderr_logfile_maxbytes=0
# autorestart=true
# startretries=3
# stopasgroup=true
# killasgroup=true
# numprocs=2

# ; sudo supervisorctl restart apidoc
# [program:apidoc]
# command=npm run watch
# # command=npm run serve
# directory=/home/ubuntu/code/demo/backend/yuqueue/Modules/doc
# stdout_logfile=/var/log/supervisor/apidoc.stdout.log
# stdout_logfile_maxbytes=0
# stderr_logfile=/var/log/supervisor/apidoc.stderr.log
# stderr_logfile_maxbytes=0
# autorestart=true
# startretries=10
# stopasgroup=true
# killasgroup=true

# ; sudo supervisorctl restart timer
# [program:timer]
# command=npm run watch
# directory=/home/ubuntu/code/demo/backend/yuqueue/Modules/timer
# stdout_logfile=/var/log/supervisor/timer.stdout.log
# stdout_logfile_maxbytes=0
# stderr_logfile=/var/log/supervisor/timer.stderr.log
# stderr_logfile_maxbytes=0
# autorestart=true
# startretries=10
# stopasgroup=true
# killasgroup=true

# ; sudo supervisorctl restart node
# [program:node]
# command=npm run watch3000
# directory=/app
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/supervisor/node.stdout.log
# stdout_logfile_maxbytes=0
# stderr_logfile=/var/log/supervisor/node.stderr.log
# stderr_logfile_maxbytes=0
# # user=ubuntu
# stopasgroup=true
# killasgroup=true


# [program:redis-server]
# # command=/etc/init.d/redis-server start
# command=/usr/bin/redis-server /etc/redis/redis.conf
# autostart=true
# autorestart=true
# user=root
# stdout_logfile=/var/log/supervisor/redis.stdout.log
# stderr_logfile=/var/log/supervisor/redis.stderr.log
# stopasgroup=true
# killasgroup=true
